# etc/supervisord.conf
[supervisord]
childlogdir=/var/log/app

# [program:nginx]
# command = /usr/sbin/nginx

[program:api]
directory = /srv
command = node build/src/index.js
autorestart = true
